# systemd/ha-light-panel.service
[Unit]
Description=HA Light Panel
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ha-light-panel
Group=ha-light-panel
WorkingDirectory=/opt/ha-light-panel
EnvironmentFile=/etc/ha-light-panel.env
ExecStart=/usr/bin/node /opt/ha-light-panel/server.js
Restart=always
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/tmp
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX

[Install]
WantedBy=multi-user.target
